# bin/unicorn-centos7_asdf.service
[Unit]
Description=SHIRASAGI Unicorn Server
After=mongod.service

[Service]
User=root
WorkingDirectory=/var/www/shirasagi
Environment=RAILS_ENV=production
SyslogIdentifier=unicorn
PIDFile=/var/www/shirasagi/tmp/pids/unicorn.pid
Type=forking
TimeoutSec=300

ExecStart=/root/.asdf/shims/bundle exec unicorn_rails -c config/unicorn.rb -D
ExecStop=/usr/bin/kill -QUIT $MAINPID
ExecReload=/usr/bin/kill -USR2 $MAINPID

[Install]
WantedBy=multi-user.target
